# semdict.service
[Unit]
Description=Semantic dict

[Service]
Type=notify
ExecStart=/usr/bin/semdict
Restart=on-failure

[Install]
WantedBy=multi-user.target
